# scripts/urlshortener-webhook.service
[Unit]
Description=URL Shortener Webhook Server
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/URL_shortner
ExecStart=/root/URL_shortner/scripts/.venv/bin/python3 /root/URL_shortner/scripts/webhook_server.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=WEBHOOK_PORT=9999
Environment=WEBHOOK_HOST=0.0.0.0

[Install]
WantedBy=multi-user.target
